# aidp-backend.service
[Unit]
Description=AIDP Backend Service
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/path/to/aidp/backend
Environment="PATH=/path/to/aidp/backend/venv/bin"
ExecStart=/path/to/aidp/backend/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 4
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aidp-backend

[Install]
WantedBy=multi-user.target
